# crontab
0 1 * * * sh /root/renewCert.sh >/dev/null 2>&1